# supervisor-config/hoogle.conf
[program:hoogle]
command=hoogle server -p 6942
