# supervisor-config/prefect.conf
[program:prefect]
command=prefect server start
